# scripts/madcatz2jstick.service
[Unit]
Description=MadCatz2Jstick

[Service]
Type=oneshot
ExecStart=/usr/local/bin/madcatz2jstick
